# templates/cliproxyapi.service
[Unit]
Description=CLIProxyAPI service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__CLIPROXY_WORKDIR__
ExecStart=__CLIPROXY_BINARY_PATH__ --config __CLIPROXY_CONFIG_PATH__
Restart=on-failure
RestartSec=3
Environment=HOME=%h

[Install]
WantedBy=default.target
